# backend/deploy/systemd/palladium-backend.service
[Unit]
Description=Palladium Backend Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=palladium
Group=palladium
WorkingDirectory=/opt/palladium/backend
EnvironmentFile=/etc/palladium/palladium-backend.env
ExecStart=/bin/bash -lc '$JAVA_HOME/bin/java -jar /opt/palladium/backend/target/palladium-backend-1.0.0.jar'
Restart=always
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/palladium/backend /var/log/palladium

[Install]
WantedBy=multi-user.target
